# playbooks/templates/nodepool-launcher/etc/systemd/system/nodepool-launcher.service.d/override.conf.j2
# This file is generated by Ansible
#  DO NOT EDIT THIS FILE BY HAND -- YOUR CHANGES WILL BE OVERWRITTEN
#
[Service]
{% if ansible_host | ipv4 %}
Environment="STATSD_HOST={{ hostvars['statsd01'].ansible_host }}"
{% endif %}
