# playbooks/templates/gear/etc/systemd/system/gear.service.d/override.conf.j2
# This file is generated by Ansible
#  DO NOT EDIT THIS FILE BY HAND -- YOUR CHANGES WILL BE OVERWRITTEN
#
[Service]
Environment="ARGS=--ssl-ca {{ gear_file_ssl_ca_dest }} --ssl-cert {{ gear_file_ssl_cert_dest }} --ssl-key {{ gear_file_ssl_key_dest }}"
Environment="PREFIX={{ gear_pip_virtualenv }}"
